# script/vpssim/menu/crontab/vpssim-befor-xoa-het-cronjob
#!/bin/sh
. /home/vpssim.conf

if [ ! -f /var/spool/cron/root ]; then
/etc/vpssim/menu/crontab/vpssim-xoa-het-cronjob
exit
fi

if [ "$(grep "auto-start-mysql" /var/spool/cron/root | awk '{print $6}')" == "auto-start-mysql" ]; then
/etc/vpssim/menu/crontab/vpssim-xoa-het-cronjob-but-auto-start-mariadb
exit
fi
if [ "$(grep "auto-start-mysql" /var/spool/cron/root | awk '{print $6}')" == "" ]; then
/etc/vpssim/menu/crontab/vpssim-xoa-het-cronjob
exit
fi
clear
echo "we can find cronjob for you"
/etc/vpssim/menu/crontab/vpssim-cronjob-menu
exit
fi
